# scripts/com.ashwell.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS launchd example. Edit the two ABSOLUTE paths below, then:

    cp scripts/com.ashwell.agent.plist ~/Library/LaunchAgents/
    launchctl load   ~/Library/LaunchAgents/com.ashwell.agent.plist
    launchctl unload ~/Library/LaunchAgents/com.ashwell.agent.plist   # to remove

  Fires weekly (Wednesday 08:00 local). Ashwell exits immediately outside the
  wake window. Weekday: 0/7=Sun, 3=Wed.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ashwell.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/ABSOLUTE/PATH/TO/ashwell/dist/cli.js</string>
    <string>run</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>3</integer>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/ashwell.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ashwell.err.log</string>
</dict>
</plist>
